# deploy/systemd/exec-clone-daily.timer
[Unit]
Description=Run Executive Clone AI daily batch at 03:00

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
