# deploy/systemd/exec-clone-weekly.timer
[Unit]
Description=Run Executive Clone AI weekly batch on Sunday 04:00

[Timer]
OnCalendar=Sun *-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
